# deploy/ansible-playbooks/templates/gunicorn.socket
[Unit]
Description=maestro socket

[Socket]
ListenStream=/run/gunicorn/gunicorn.sock

[Install]
WantedBy=sockets.target